# scripts/tribunal-loop.service
[Unit]
Description=Tribunal Quota-Aware Loop (gu-log)
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/clawd/projects/gu-log
ExecStart=/bin/bash %h/clawd/projects/gu-log/scripts/cc-tribunal-loop-wrapper.sh --workers 5
Restart=on-failure
RestartSec=60
Environment=TZ=Asia/Taipei

# Graceful stop (see OpenSpec: tribunal-graceful-run-control)
# KillMode=mixed — SIGTERM goes to main process (wrapper→loop) only; the
# tribunal-all-claude.sh subprocess is NOT killed and is allowed to finish
# the current article. The loop's trap sets stop_requested and exits
# cleanly at the next safe boundary.
KillMode=mixed
KillSignal=SIGTERM
# Worst case: article is in the middle of rewrite + build + frontmatter
# write. Allow up to 60min before SIGKILL.
TimeoutStopSec=3600

# Resource limits — tuned for --workers 5 (5 parallel tribunal-all-claude.sh
# subprocesses, each spawning a claude CLI judge + optional pnpm build).
# CPUQuota 200% lets bursts of concurrent builds use multiple cores without
# starving each other; MemoryMax 2G covers peak (5 × ~150MB worker +
# supervisor + shared helpers) with headroom for pnpm node_modules caching.
CPUQuota=200%
MemoryMax=2G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tribunal-loop

[Install]
WantedBy=default.target
